# service/lenovo-fancurve-restart.service
[Service]
Type=oneshot
ExecStart=/bin/systemctl restart lenovo-fancurve.service

[Unit]
StartLimitIntervalSec=0