# service/lenovo-fancurve-restart.path
[Path]
Unit=lenovo-fancurve-restart.service
PathModified=/sys/firmware/acpi/platform_profile
PathModified=/sys/class/power_supply/ADP0/online

[Install]
WantedBy=default.target
