# service/lenovo-fancurve.save
#!/sbin/openrc-run

extra_started_commands="reload"
command="/etc/lenovo-fan-control/fancurve-set.sh"
description="Apply Fan Curve to Lenovo Legion Laptop depending of Power State and Charger Mode (LenovoLegionLinux kernel Module)"

depend() {
	after agetty
	use logger
}

reload() {
	ebegin "Reloading lenovo-fancurve"
	start-stop-daemon --exec $command --signal HUP
	eend $?
}
