# service/lenovo-fancurve-restart.service
[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart lenovo-fancurve.service